# infra/systemd/tronsoftos.service
[Unit]
Description=TronSoftOS host service
After=network-online.target firebird.service
Wants=network-online.target

[Service]
Type=simple
User=tronsoftos
Group=tronsoftos
WorkingDirectory=/opt/tronsoftos
EnvironmentFile=-/etc/tronsoftos/tronsoftos.env
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p /opt/tronsoftos/state /opt/tronsoftos/logs /opt/tronsoftos/config/rclone
ExecStartPre=/bin/chown -R tronsoftos:tronsoftos /opt/tronsoftos/state /opt/tronsoftos/logs /opt/tronsoftos/config/rclone
ExecStartPre=/bin/chmod 700 /opt/tronsoftos/state
ExecStartPre=/bin/chmod 700 /opt/tronsoftos/config/rclone
ExecStart=/usr/bin/node /opt/tronsoftos/backend/src/server.mjs
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
